# capture/systemd/skylapse-capture.service
[Unit]
Description=Skylapse Timelapse Capture Service
Documentation=https://github.com/skylapse/docs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=exec
User=skylapse
Group=video
WorkingDirectory=/opt/skylapse/capture

# Python environment
Environment=PYTHONPATH=/opt/skylapse/capture/src
Environment=PYTHONUNBUFFERED=1

# Service configuration
Environment=SKYLAPSE_ENV=production
Environment=SKYLAPSE_LOG_LEVEL=INFO
Environment=SKYLAPSE_CONFIG_DIR=/opt/skylapse/capture/config

# Performance settings
Environment=OMP_NUM_THREADS=1

# Execute the service
ExecStart=/opt/skylapse/capture/venv/bin/python -m src.capture_service
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=10

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/skylapse
SupplementaryGroups=gpio i2c spi

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skylapse-capture

# Watchdog - disabled for now since service doesn't implement watchdog notifications
# WatchdogSec=30

[Install]
WantedBy=multi-user.target
Alias=skylapse.service
